# mysql_backup/mysql_backup.cron
*/1 * * * * sh  ~/shell/mysql_backup.sh
